# metal3-locallab/com.dave.ironic-portforward.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dave.ironic-portforward</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/kubectl</string>
        <string>port-forward</string>
        <string>-n</string>
        <string>metal3-system</string>
        <string>svc/metal3-metal3-ironic</string>
        <string>6185:6185</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/ironic-portforward.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ironic-portforward.log</string>
</dict>
</plist>

